# crontab
14 14 14 2 5 /bin/bash /home/absolutlubis/Documents/soalshift1/no1.sh

2-30/6 * * * * /bin/bash /home/absolutlubis/Documents/soalshift1/no5.sh
